# netif.d/bridge
#!/sbin/runscript 

depend() {
	config "$settings"
	[ -n "$slaves" ] && need $slaves
}

start() {
	ebegin "$STARTD"

	# First, create bridge, add optional slave interfaces:

	brctl addbr $interface || die "Couldn't create bridge $interface"
	if [ -n "$slaves" ]
	then
		local slave
		for slave in ${slaves//netif./}
		do
			brctl addif $interface $slave || die "Couldn't add slave $slave to bridge $interface"
		done
	fi

	# Set forwarding delay to 0 - good for non-spanning-tree/virtualization use:
	
	brctl setfd $interface 0 || die "Couldn't set forwarding delay to zero on bridge $interface"

	# Configure interface:
	
	ip link set $interface up || die "Couldn't bring $interface up"
	[ -n "$mtu" ] && { ip link set dev $interface mtu $mtu || die "Couldn't set MTU on $interface to $mtu"; }
        ezroute add
	ezresolv add
	eend 0
}

stop() {
	ebegin "$STOPD"
	ezresolv del
	ezroute del
	ip link set $interface down || die "Couldn't bring $interface down"
	brctl delbr $interface
	eend 0
}
