# netif.d/bridge-noip
#!/sbin/runscript

depend() {
	config "$CONFD"
	[ -n "$slaves" ] && need $slaves
}

start() {
	local slave
	ebegin "$STARTD"
	brctl addbr $interface
	if [ -n "$slaves" ]
	then
		for slave in ${slaves//netif./}
		do
			brctl addif $interface ${slave}
		done
	fi	
	brctl setfd $interface 0 && \
	ip link set $interface up && \
	eend $?
}

stop() {
	ebegin "$STOPD"
	ip link set $interface down && \
	brctl delbr $interface
	eend $?
}
