# init.d/netif.lo
#!/sbin/runscript
# Copyright (c) 2007-2009 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

start() {
        ebegin "Bringing up network interface lo"
        ifconfig lo 127.0.0.1 netmask 255.0.0.0 up && \
        route add -net 127.0.0.0 netmask 255.0.0.0 dev lo
        eend $?
}

stop() {
        ebegin "Shutting down network interface lo"
        route del -net 127.0.0.0 netmask 255.0.0.0 dev lo &&
        ifconfig lo down
        eend $?
}
